# supervisor-bot.conf
[program:trading-bot]
directory=/home/ubuntu
command=python3 run_bot.py
user=ubuntu
autostart=true
autorestart=false
stderr_logfile=/home/ubuntu/bot-err.log
stdout_logfile=/home/ubuntu/bot-log.log